# templates/default/upstart.conf.erb
description "zipkin-<%= @service %>"
author "Rich Schumacher <rich@curalate.com>"

start on stopped cloud-final
stop on shutdown

respawn
respawn limit 3 60

setuid <%= @user %>
setgid <%= @user %>

script
    cd <%= @path %>

    CLASSPATH=""
    for JAR in <%= @path %>/build/libs/*.jar ; do
        CLASSPATH=$CLASSPATH:$JAR
    done

    <% @env_vars.each do |k, v| %>
    <%= "export #{k}=\"#{v}\"" %>
    <% end %>

    java -cp ${CLASSPATH#?} com.twitter.zipkin.<%= @service %>.Main <%= @jvm_args ? @jvm_args.map { |k, v| ["-#{k}=#{v}"] }.join(' ') : '' %> <%= @service != 'web' ? "-f config/#{@service}-#{@config_key}.scala" : '' %> >> /var/log/zipkin-<%= @service %>.log 2>&1
end script
