# templates/default/upstart.conf.erb
description "zipkin"
author "Rich Schumacher <rich@curalate.com>"

start on stopped cloud-final
stop on shutdown

respawn
respawn limit 3 60

setuid <%= @user %>
setgid <%= @user %>

script
    <% zipkin_env_vars.each do |k, v| %>
    <%= "export #{k}=\"#{v}\"" %>
    <% end %>

    chdir <%= zipkin_version_dir %>
    java \
        -Dloader.path='<%= zipkin_kafka_jar_file %>,<%= zipkin_kafka_jar_file %>!/lib' \
        -Dspring.profiles.active=kafka \
        -cp <%= zipkin_jar_file %> \
        org.springframework.boot.loader.PropertiesLauncher
end script
